# supervisord.conf
[supervisord]
nodaemon=true
[program:shadowsocks]
command=/etc/init.d/shadowsocks -D
[program:finalspees]
command=/bin/bash /fs/fs_start.sh
